# homepage.service
[Unit]
Description=bertolsson.com

[Service]
ExecStart=/var/www/html/home/homestart.sh

[Install]
WantedBy=multi-user.target
