# homepage.service
[Unit]
Description=dev.bertolsson.com

[Service]
ExecStart=/var/www/html/home-dev/homestart.sh

[Install]
WantedBy=multi-user.target
